# host/etc/systemd/system/restic-backup-appdata.timer
[Unit]
Description=Nightly restic backup of /docker/appdata

[Timer]
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=30min
Persistent=true

[Install]
WantedBy=timers.target
