# host/etc/systemd/system/restic-backup-appdata.service
[Unit]
Description=Restic backup of /docker/appdata
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/restic-backup-appdata.sh
Nice=10
IOSchedulingClass=idle
